# ICYOU.Server.Linux/icyou-server.service
[Unit]
Description=ICYOU Server - Мессенджер сервер
After=network.target

[Service]
Type=simple
User=psyshout
Group=psyshout
WorkingDirectory=/home/psyshout/publish_linux
ExecStart=/home/psyshout/publish_linux/run_server.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=icyou-server

# Переменные окружения
Environment="DOTNET_ROOT=/usr/share/dotnet"
Environment="ASPNETCORE_ENVIRONMENT=Production"

# Безопасность
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

